# deploy/s6/open-webui-livekit-agent/run
#!/bin/sh
set -eu

cd /var/www/open-webui

export HOME=/home/arch
export TMPDIR=/home/arch/tmp/open-webui-livekit-agent-runtime

set -a
. /var/www/open-webui/.env
if [ -f /var/www/open-webui/deploy/livekit/livekit.env ]; then
  . /var/www/open-webui/deploy/livekit/livekit.env
fi
set +a

export LIVEKIT_URL="ws://127.0.0.1:7880"

exec /usr/bin/s6-setuidgid arch \
  /var/www/open-webui/.venv/bin/python /var/www/open-webui/deploy/livekit/agent.py start
