# deploy/s6/open-webui-livekit-server/run
#!/bin/sh
set -eu

cd /home/pepper/apps/open-webui

export HOME=/home/pepper
export TMPDIR=/home/pepper/tmp/open-webui-livekit-server-runtime

# Base env (provider keys, etc.)
set -a
. /home/pepper/apps/open-webui/.env

# LiveKit URL + credentials (generated by deploy/livekit/setup_self_hosted.py)
if [ -f /home/pepper/apps/open-webui/deploy/livekit/livekit.env ]; then
  . /home/pepper/apps/open-webui/deploy/livekit/livekit.env
fi
set +a

exec /usr/bin/s6-setuidgid pepper \
  /usr/local/bin/livekit-server --config /home/pepper/apps/open-webui/deploy/livekit/livekit.yaml

