# deploy/s6/open-webui-livekit-portal/run
#!/bin/sh
set -eu

cd /var/www/open-webui

export HOME=/home/arch
export TMPDIR=/home/arch/tmp/open-webui-livekit-portal-runtime
export PORT=8091

set -a
. /var/www/open-webui/.env
if [ -f /var/www/open-webui/deploy/livekit/livekit.env ]; then
  . /var/www/open-webui/deploy/livekit/livekit.env
fi
set +a

exec /usr/bin/s6-setuidgid pepper \
  /var/www/open-webui/.venv/bin/python /var/www/open-webui/deploy/livekit/portal.py
