# deploy/s6/open-webui/run
#!/bin/sh
set -eu

cd /home/pepper/apps/open-webui

export HOME=/home/pepper
export TMPDIR=/home/pepper/tmp/open-webui-runtime

set -a
. /home/pepper/apps/open-webui/.env
set +a

# Sync OpenAI-compatible provider keys from env → Postgres so updates to ~/.api-keys
# take effect without needing to click through the admin UI.
/usr/bin/s6-setuidgid pepper \
  /home/pepper/apps/open-webui/.venv/bin/python /home/pepper/apps/open-webui/deploy/sync_openai_provider_keys_from_env.py || true

# Ensure provider prefixes exist (prevents de-duplication across providers).
/usr/bin/s6-setuidgid pepper \
  /home/pepper/apps/open-webui/.venv/bin/python /home/pepper/apps/open-webui/deploy/ensure_openai_provider_prefixes.py || true

exec /usr/bin/s6-setuidgid pepper \
  /home/pepper/apps/open-webui/.venv/bin/open-webui serve --host 127.0.0.1 --port 8080
